# leetcode-tracker.service
[Unit]
Description=LeetCode Tracker FastAPI Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/leetcode_tracker_uv
Environment="PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/root/leetcode_tracker_uv/.env
ExecStart=/root/.local/bin/uv run uvicorn leetcode_tracker.main:app --host 127.0.0.1 --port 8000 --proxy-headers
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
